# systemd/ics-backup.service
[Unit]
Description=ICS_DMAS 指揮部 SQLite 自動備份（每日 03:00 by ics-backup.timer）
After=ics-command.service
ConditionPathExists=/home/ics/ics-dmas/command-dashboard/data/ics.db

[Service]
Type=oneshot
User=ics
WorkingDirectory=/home/ics/ics-dmas/command-dashboard
ExecStart=/home/ics/ics-dmas/command-dashboard/.venv/bin/python /home/ics/ics-dmas/command-dashboard/scripts/backup_db.py
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
StandardOutput=journal
StandardError=journal
